# system_files/etc/profile.d/botany-alias.sh
alias botany_sudo='/usr/sbin/sudo -u botany_adm /usr/sbin/sudo'
alias botany_su='/usr/sbin/sudo -u botany_adm /usr/sbin/sudo /usr/sbin/su -'
alias botany_adm='/usr/sbin/sudo -u botany_adm -i ;'

alias bsudo=botany_sudo
alias bsu=botany_su
alias badm=botany_adm

alias edit='/usr/bin/msedit'
