# system_files/usr/lib/systemd/system/beesd@.service.d/00-bees-overrides.conf
[Unit]
ConditionACPower=true

[Service]
ExecStart=
ExecStart=/usr/bin/beesd --no-timestamps --strip-paths --thread-count=1 --throttle-factor=8 %i
#ExecStart=/usr/bin/timeout 15m /usr/bin/beesd --no-timestamps --strip-paths --thread-count=1 --throttle-factor=10 %i
#SuccessExitStatus=124
Restart=no
OOMScoreAdjust=1000
OOMPolicy=stop
#CPUQuota=50%
