# system_files/usr/lib/systemd/system/beesd@.service.d/00-bees-overrides.conf
[Service]
ExecStart=
ExecStart=/usr/bin/timeout 15m /usr/bin/beesd --no-timestamps --strip-paths --thread-count=1 --throttle-factor=10 %i
Restart=no
SuccessExitStatus=124
OOMScoreAdjust=1000
OOMPolicy=stop
